# cron_jobs/weekly_backup.cron
0 2 * * 0 /home/ubuntu/network_automation/scripts/run_backup.sh >> /home/ubuntu/network_automation/logs/backup_cron.log 2>&1
